# sensor.deb/DEBIAN/prerm
#!/bin/bash
/etc/init.d/sensor stop
/usr/sbin/update-rc.d sensor disable
/usr/sbin/update-rc.d sensor remove